# docker/s6/gunicorn/run
#!/command/with-contenv /command/execlineb
cd /app/empi-server
poetry run gunicorn empi_server.wsgi --bind 127.0.0.1:8001 --workers 4 --max-requests 1000
