# docker/s6/node/run
#!/command/with-contenv /command/execlineb
cd /app/empi-web
node .
