# docker/s6/migrations/up
/command/with-contenv sh -c "cd /app/empi-server; poetry run ./manage.py migrate"
